# pipewire/.config/service/pipewire/run
#!/bin/sh

[ -r ./conf ] && . ./conf

exec chpst -e "$TURNSTILE_ENV_DIR" pipewire $OPTS 2>&1
